# supervisord.conf
[supervisord]
nodaemon=true

[program:apache]
command=/usr/sbin/apache2ctl -D FOREGROUND

[program:cron]
command=cron -f
autostart=true
autorestart=true